# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=service httpd restart
autorestart=false

[program:nagios]
command=service nagios restart
autorestart=false
